# conf/etc/s6/services/02-nginx/run
#!/bin/sh
NAME="NGINX" 
COMMAND="nginx -g 'daemon off;'" 
WAIT="/tmp/kanboard"

exec service-handler "$NAME" "$COMMAND" "$WAIT"
